# cron/crontab
#crontab
@reboot root *data_feed_tasks/sin_wave_over_time_in_seconds